# rootfs/etc/s6-overlay/s6-rc.d/ssh-agent/run
#!/bin/bash
set -eao pipefail
source /etc/ssh/sshagentrc
echo "ssh-agent: INFO: Starting ssh-agent service..."
ssh-agent -D -a ${SSH_AUTH_SOCK} >/dev/null
